# tools/mlrift-uefi-probe.service
[Unit]
Description=MLRift η.3f UEFI POST state probe (read-only)
Documentation=https://github.com/Rift-Intelligence/MLRift examples/phase3_eta3f_post_uefi_state.mlr
# Run after vfio-pci has had a chance to bind to the dGPU. multi-user
# target reaches AFTER kernel modules + udev settle.
After=multi-user.target systemd-modules-load.service
Wants=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Working directory matters because the smoke uses relative paths in
# any imports the runtime might resolve. Project root.
WorkingDirectory=/home/pantelis/Desktop/Projects/Work/MLRift

# Log header for correlation across boots.
ExecStartPre=/bin/sh -c 'echo "=== MLRift UEFI POST probe — boot at $(date -Iseconds) ===" >> /var/log/mlrift_uefi_probe.log'

# Actual smoke. Output streamed to per-boot log marker.
ExecStart=/home/pantelis/Desktop/Projects/Work/MLRift/build/phase3_eta3f_post_uefi_state

# Footer + chmod so user can read without sudo afterwards.
ExecStartPost=/bin/sh -c 'echo "=== probe complete ===" >> /var/log/mlrift_uefi_probe.log; chmod 644 /var/log/mlrift_uefi_probe.log'

StandardOutput=append:/var/log/mlrift_uefi_probe.log
StandardError=append:/var/log/mlrift_uefi_probe.log

# Read-only smoke. No need for restart on failure.
Restart=no

# Root needed for /dev/vfio access.
User=root

[Install]
WantedBy=multi-user.target
